# phosphor-pid-autotune.service
[Unit]
Description=Phosphor PID Autotune
Wants=xyz.openbmc_project.EntityManager.service
After=dbus.service xyz.openbmc_project.EntityManager.service

[Service]
Type=simple
# Safety guard: ensure the binary exists before starting
ExecStartPre=/usr/bin/test -x /usr/bin/phosphor-pid-autotune
# Pass an explicit config path (falls back to EntityManager first inside the app)
ExecStart=/usr/bin/phosphor-pid-autotune /usr/share/phosphor-pid-autotune/configs/autotune.json
Restart=on-failure
RestartSec=3
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phosphor-pid-autotune

[Install]
WantedBy=multi-user.target
